# system/udev/99-mycolab-usb-backup.rules
ACTION=="add", KERNEL=="sd?1", SUBSYSTEM=="block", RUN+="/usr/bin/systemd-run /home/zerocool/My-co-Lab/scripts/usb_backup.sh %k"
